# service/rc.gentoo
#!/sbin/runscript

depend()
{
	need localmount net
	after bootmisc modules
}

start()
{
	
	ebegin "Starting EmuFox virtual machines"
	/usr/bin/emufox-autostart
	eend $?
}

stop()
{
	ebegin "Stopping EmuFox virtual machines"
	/usr/bin/emufox-autostop
	eend $?
}
